# data/90-ppm-brightness.rules
ACTION=="add", SUBSYSTEM=="backlight", RUN+="ppm-bright-helper video g+w /sys/class/backlight/%k/brightness"
